# shopply.conf
description "Tornado Service"
author  "Salil Wadnerkar"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/ubuntu/helloshopply
exec ./hello_shopply.py
respawn
